# systemd/sim7600-lte.service
[Unit]
Description=SIM7600 LTE Data Connection (ECM/USB Ethernet)
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/home/daniel/rpi-datalogger/.venv/bin/python3 /opt/sim7600/lte-monitor.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
